# rel/overlays/etc/backup.cron
00 17 * * * /rails/bin/litestream-backup
